# larch7/larch/profiles/xmini/rootoverlay/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

x10d_fstab

# Restore saved sound volume, etc.
alsactl restore
